# gpu_monitor/deploy/celery-ingest@.service
[Unit]
Description=Celery Ingest Worker %i
After=network.target redis.service postgresql.service
Wants=redis.service postgresql.service

[Service]
Type=simple
User=monitoring
Group=monitoring
WorkingDirectory=/opt/gpu_monitor
EnvironmentFile=/opt/gpu_monitor/.env
ExecStart=/opt/gpu_monitor/venv/bin/celery -A gpu_monitor worker \
    --loglevel=INFO \
    --queues=ingest \
    --concurrency=2 \
    --pool=prefork \
    --hostname=ingest-worker-%i@%h \
    --max-tasks-per-child=100 \
    --time-limit=300 \
    --soft-time-limit=240
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
